# Log/supervisord.conf
[supervisord]
nodaemon=true

[program:bash]
command=/bin/bash -c "/bin/bash"
user:user

# [program:Log]
# command=/bin/bash -c "python /user-script/Log/Log.py"
# suser:user
# username=root
# password=docke
